# crontab
# cronjob file
*/2 * * * * /usr/bin/python /home/python/PyMail-MailJet.py